# packer/ansible/roles/supervisord/templates/supervisord.conf

[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/bin/php-fpm

[program:nginx]
command=/usr/sbin/nginx

[program:wordpress]
command=/install_wp.sh